# deploy/temperature_daemon.conf
# Database updater

description "Temperature Updater"
start on runlevel [2345]
stop on runlevel [06]

respawn

exec /var/www/alice_pi/env/bin/python /var/www/alice_pi/bin/temperature_daemon
